# server/systemd/udtun-server.service
[Unit]
Description=UDTUN UDP Tunneling Server
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/udtun/server
ExecStart=/usr/bin/python3 /opt/udtun/server/main.py /etc/udtun/server.json
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=udtun-server

# Security
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/log/udtun

# Performance
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
